# system/systemd/flugradar-display.service
# Template — install.sh generates the actual service file with correct paths.
# This file serves as reference only.

[Unit]
Description=Sasso Radar Tower — Display
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=%u
EnvironmentFile=-/home/%u/.env
ExecStart=/home/%u/sasso-radar-tower/system/flugradar-display-start.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
