# system/systemd/flugradar-web.service
# Template — install.sh generates the actual service file with correct paths.
# This file serves as reference only.

[Unit]
Description=Sasso Radar Tower — Web Portal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%u
EnvironmentFile=-/home/%u/.env
ExecStart=/home/%u/sasso-radar-tower/.venv/bin/flugradar-web --host 0.0.0.0 --port 5000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
